# rpm-builder/rippled.service
[Unit]
Description=Ripple Daemon

[Service]
Type=simple
Environment="LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/opt/ripple/openssl/lib:/opt/ripple/boost/lib"
ExecStart=/opt/ripple/bin/wrapper.sh --conf /opt/ripple/etc/rippled.cfg
ExecStop=/opt/ripple/bin/rippled --conf /opt/ripple/etc/rippled.cfg stop
Restart=no
User=rippled
Group=rippled

[Install]
WantedBy=multi-user.target
